# templates/default/initd.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:          solr
# Required-Start:   
# Required-Stop:
# Default-Start:     3 5                                                                                                               
# Default-Stop:      0 1 2 6
# Short-Description: init script for solr search                                                                                       
# Description:       init script for solr search                                                                                       
### END INIT INFO                                                                                                                      
                                                                                                                                       
# source function library                                                                                                              
. /lib/lsb/init-functions                                                                                                              

PIDFILE=/var/run/solr.pid                                                                                                              
                                                                                                                                       
LOG_FILE=/var/log/solr.log

NAME="Solr"                                                                                                                            
        
DESC="start/stop Solr Server"                                                                                                          
                                                                                                                                       
DAEMON="/var/lib/solr.start"

start() {
    if start-stop-daemon --stop --quiet --pidfile $PIDFILE --signal 0                                                                  
    then                                                                                                                               
        echo $NAME was already running: $PID                                                                                           
        exit 2;                                                                                                                        
    else
        start-stop-daemon --start --pidfile $PIDFILE --exec $DAEMON                                                                    
    fi
}                                                                                                                                      
        
stop() {
        start-stop-daemon --stop --retry=TERM/30/KILL/5 --pidfile $PIDFILE                                                             
}       
        
case $1 in      
    start)
        start   
        ;;      
    stop)                                                                                                                              
        stop                                                                                                                           
        ;;
    restart)
        stop                                                                                                                           
        sleep 2                                                                                                                        
        start                                                                                                                          
        ;;
    status)                                                                                                                            
        status_of_proc -p $PIDFILE "$DAEMON" "$NAME" && exit 0 || exit $?                                                              
        ;;                                                                                                                             
    *)  
        echo "Usage: $0 {start|stop|restart|status}"                                                                                   
        exit 3
        ;;
esac
